# files/rtcwake-for-mc-server-service-any@any/rtcwake-for-mc-server.service
[Unit]
Description=Suspend at 01:30+UTC and wake at 07:00+UTC.

[Service]
Type=oneshot
ExecStart=/usr/local/bin/do-rtcwake-for-mc-server
